# com.claude-space.launcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-space.launcher</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-l</string>
        <string>-c</string>
        <!-- Update this path to your actual location -->
        <string>$HOME/claude-space/claude-space-launcher.sh</string>
    </array>
    
    <!-- Run every 2 hours (7200 seconds) -->
    <key>StartInterval</key>
    <integer>7200</integer>
    
    <!-- Also run when first loaded -->
    <key>RunAtLoad</key>
    <false/>
    
    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/tmp/claude-space-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-space-stderr.log</string>
    
    <!-- Environment variables (optional overrides) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>CLAUDE_SPACE_DIR</key>
        <string>$HOME/claude-space</string>
        <key>SESSION_DURATION_SECONDS</key>
        <string>900</string>
        <key>QUOTA_THRESHOLD</key>
        <string>30</string>
    </dict>
    
    <!-- Keep alive settings (restart if crashes) -->
    <key>KeepAlive</key>
    <false/>
    
    <!-- Nice value (lower priority) -->
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>

<!--
INSTALLATION:

1. Copy this file to ~/Library/LaunchAgents/:
   cp com.claude-space.launcher.plist ~/Library/LaunchAgents/

2. Edit the file to fix the $HOME path (launchd doesn't expand $HOME):
   sed -i '' "s|\$HOME|$HOME|g" ~/Library/LaunchAgents/com.claude-space.launcher.plist

3. Load the agent:
   launchctl load ~/Library/LaunchAgents/com.claude-space.launcher.plist

4. Check status:
   launchctl list | grep claude-space

5. To unload:
   launchctl unload ~/Library/LaunchAgents/com.claude-space.launcher.plist

6. To run manually (test):
   launchctl start com.claude-space.launcher

7. View logs:
   tail -f /tmp/claude-space-stdout.log
   tail -f /tmp/claude-space-stderr.log
-->
